# ansible/roles/axed-generate-miner/templates/axed-generate-miner.service.j2
[Unit]
After=network.target

[Service]
User={{ axed_user }}
Group={{ axed_group }}
ExecStart=/usr/local/bin/axed-generate-miner.sh

[Install]
WantedBy=default.target
